# .crontab
# m h  dom mon dow   command
58 16 * * FRI bash -c 'DISPLAY=:0 XDG_RUNTIME_DIR=/run/user/1000 notify-send "Cron run $(date -d "Friday 17:00" "+\%A \%H:\%M")" "The week no. $(date +\%W)  ends now" -u CRITICAL'
00 17 * * FRI  bash -c "DISPLAY=:0 XDG_RUNTIME_DIR=/run/user/1000 /usr/local/bin/cvlc --fullscreen --play-and-exit --no-loop --no-repeat --no-keyboard-events --no-mouse-events /home/b/Videos/billie_eilish.mp4"  # https://www.youtube.com/watch?v=I4G9MZ_SQmo

